# pkg/aur/praxis/praxis-service.service
[Unit]
Description=Praxis Service
PartOf=praxis.service
After=rabbitmq.service

[Service]
Type=simple
ExecStart=/usr/bin/praxis_service
EnvironmentFile=/etc/praxis/env
User=praxis
Group=praxis
Restart=on-failure
RestartSec=3

[Install]
WantedBy=praxis.service
