# pkg/systemd/praxis-web.service
[Unit]
Description=Praxis Web
After=network.target praxis-service.service
Wants=praxis-service.service

[Service]
Type=simple
ExecStart=/usr/local/bin/praxis_web
EnvironmentFile=-/etc/praxis/env
User=praxis
Group=praxis
WorkingDirectory=/var/lib/praxis
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
